# usb-watcher.service
[Unit]
After=network.target

[Service]
ExecStart=

[Install]
WantedBy=default.target
